# image/templates/files/compliance-postboot.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
    Copyright 2025 Oxide Computer Company
-->

<service_bundle type='manifest' name='postboot'>

<service name='site/postboot' type='service' version='1'>
    <create_default_instance enabled='true' />

    <!-- Wait for multi-user... -->
    <dependency name='multi-user-server' grouping='require_all'
        restart_on='none' type='service'>
        <service_fmri value='svc:/milestone/multi-user-server' />
    </dependency>

    <!-- ... and make sure we run after the T6 has been configured and the
         cxgbe links have been assigned link-local v6 addresses. -->
    <dependency name='net-setup-cxgbe' grouping='require_all' restart_on='none'
        type='service'>
        <service_fmri value='svc:/oxide/net-setup:cxgbe' />
    </dependency>

    <!-- Compliance postboot is just a milestone service. -->
    <exec_method type='method' name='start' exec=':true' timeout_seconds='30' />
    <exec_method type='method' name='stop' exec=':true' timeout_seconds='30' />

    <property_group name='startd' type='framework'>
        <propval name='duration' type='astring' value='transient' />
    </property_group>

    <stability value='Unstable' />

    <template>
        <common_name>
            <loctext xml:lang='C'>compliance postboot</loctext>
        </common_name>
        <description>
            <loctext xml:lang='C'>milestone indicating we started successfully
            </loctext>
        </description>
    </template>
</service>

</service_bundle>
