# image/templates/files/mfg.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
    Copyright 2024 Oxide Computer Company
-->

<service_bundle type="manifest" name="site-mfg">
    <service name="site/mfg" type="service" version="1">
        <create_default_instance enabled="true" />

        <dependency name="multi-user" grouping="require_all" restart_on="error"
            type="service">
            <service_fmri value="svc:/milestone/multi-user:default" />
        </dependency>

        <!-- Make sure we have an appropriate NIC setup: -->
        <dependency name='net' grouping='require_any' restart_on='none'
            type='service'>
            <service_fmri value='svc:/oxide/net-setup:igb' />
            <service_fmri value='svc:/oxide/net-setup:e1000g' />
        </dependency>

        <exec_method type="method" name="start" timeout_seconds="30"
            exec="/lib/svc/method/mfg" />

        <exec_method type="method" name="stop" exec=":kill"
            timeout_seconds="30" />

        <property_group name="startd" type="framework">
            <propval name="duration" type="astring" value="child" />
        </property_group>

        <template>
            <common_name>
                <loctext xml:lang="C">manufacturing image network setup
                </loctext>
            </common_name>
        </template>
    </service>
</service_bundle>
